# assets/99-ftdi-latency-timer.rules
ACTION=="add", SUBSYSTEMS=="tty", ENV{ID_VENDOR_ID}=="0403", RUN+="__TARGET_DIR__/ftdi-latency-timer.sh  %S$devpath"
